# recipes-tools/startupscript/files/gpio.service
[Unit]
Description=gpio service

[Service]
Type=oneshot
User=root
ExecStart=/bin/sh /etc/gpio_enable.sh 

[Install]
WantedBy=multi-user.target
