# recipes-tools/startupscript/files/gui.service
[Unit]
Description=gui service

[Service]
Type=oneshot
User=root
ExecStart=/bin/sh /etc/gui.sh

[Install]
WantedBy=multi-user.target
